# services/payment-svc/systemd/payment-svc.service
[Unit]
Description=Payment service (AIOps demo)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/payment-svc
Environment=KAFKA_BOOTSTRAP=localhost:9092
EnvironmentFile=-/etc/payment-svc.env
ExecStart=/opt/payment-svc/venv/bin/gunicorn -b 0.0.0.0:5000 -w 2 app:app
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
